# addon/rootfs/etc/s6-overlay/s6-rc.d/pi-service/finish
#!/command/with-contenv bashio
# shellcheck shell=bash
bashio::log.warning "Pi Agent HTTP service stopped (exit code: ${1})"
