# addon/rootfs/etc/s6-overlay/s6-rc.d/init-pi/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Pi Agent Add-on: Environment setup
# Sets up PI_CODING_AGENT_DIR and environment variables
# ==============================================================================

bashio::log.info 'Setting up Pi Agent environment...'

# --- Persistent storage ---
PI_DIR=/data/pi-agent
mkdir -p "${PI_DIR}/sessions"

# --- Docs cache ---
mkdir -p /data/ha-docs

# --- Write s6 env vars (available to all services via with-contenv) ---
printf '%s' "${PI_DIR}" > /var/run/s6/container_environment/PI_CODING_AGENT_DIR
printf '%s' "${SUPERVISOR_TOKEN}" > /var/run/s6/container_environment/HA_TOKEN
printf '%s' "http://supervisor/core" > /var/run/s6/container_environment/HA_URL
printf '%s' "/homeassistant" > /var/run/s6/container_environment/HA_CONFIG_PATH
printf '%s' "/data/ha-docs" > /var/run/s6/container_environment/HA_DOCS_DATA_DIR

# --- Write profile.d for interactive shells (tmux/ttyd) ---
{
    echo "export PI_CODING_AGENT_DIR=\"${PI_DIR}\""
    echo "export SUPERVISOR_TOKEN=\"${SUPERVISOR_TOKEN}\""
    echo "export HA_TOKEN=\"${SUPERVISOR_TOKEN}\""
    echo "export HA_URL=\"http://supervisor/core\""
    echo "export HA_CONFIG_PATH=\"/homeassistant\""
    echo "export HA_DOCS_DATA_DIR=\"/data/ha-docs\""
} > /etc/profile.d/pi-agent.sh

# --- Default provider ---
default_provider=$(bashio::config 'default_provider')
if bashio::var.has_value "${default_provider}"; then
    printf '%s' "${default_provider}" > /var/run/s6/container_environment/PI_DEFAULT_PROVIDER
    echo "export PI_DEFAULT_PROVIDER=\"${default_provider}\"" >> /etc/profile.d/pi-agent.sh
fi

# --- Default model ---
default_model=$(bashio::config 'default_model')
if bashio::var.has_value "${default_model}"; then
    printf '%s' "${default_model}" > /var/run/s6/container_environment/PI_DEFAULT_MODEL
    echo "export PI_DEFAULT_MODEL=\"${default_model}\"" >> /etc/profile.d/pi-agent.sh
fi

# --- Service provider/model (for pi_agent.ask, separate from interactive defaults) ---
service_provider=$(bashio::config 'service_provider')
if bashio::var.has_value "${service_provider}"; then
    printf '%s' "${service_provider}" > /var/run/s6/container_environment/PI_SERVICE_PROVIDER
    echo "export PI_SERVICE_PROVIDER=\"${service_provider}\"" >> /etc/profile.d/pi-agent.sh
fi

service_model=$(bashio::config 'service_model')
if bashio::var.has_value "${service_model}"; then
    printf '%s' "${service_model}" > /var/run/s6/container_environment/PI_SERVICE_MODEL
    echo "export PI_SERVICE_MODEL=\"${service_model}\"" >> /etc/profile.d/pi-agent.sh
fi

# --- API keys from dedicated config fields ---
declare -A API_KEY_MAP=(
    [anthropic_api_key]=ANTHROPIC_API_KEY
    [openai_api_key]=OPENAI_API_KEY
    [google_api_key]=GEMINI_API_KEY
    [openrouter_api_key]=OPENROUTER_API_KEY
    [xai_api_key]=XAI_API_KEY
    [groq_api_key]=GROQ_API_KEY
    [mistral_api_key]=MISTRAL_API_KEY
    [cerebras_api_key]=CEREBRAS_API_KEY
    [huggingface_token]=HF_TOKEN
    [github_token]=GITHUB_TOKEN
    [aws_access_key_id]=AWS_ACCESS_KEY_ID
    [aws_secret_access_key]=AWS_SECRET_ACCESS_KEY
    [aws_region]=AWS_REGION
    [azure_openai_api_key]=AZURE_OPENAI_API_KEY
)

for config_key in "${!API_KEY_MAP[@]}"; do
    env_var="${API_KEY_MAP[$config_key]}"
    value=$(bashio::config "${config_key}")
    if bashio::var.has_value "${value}"; then
        printf '%s' "${value}" > "/var/run/s6/container_environment/${env_var}"
        echo "export ${env_var}=\"${value}\"" >> /etc/profile.d/pi-agent.sh
        bashio::log.info "Set ${env_var} from config"
    fi
done

# --- Persistent bash history ---
touch /data/.bash_history
ln -sf /data/.bash_history /root/.bash_history

# --- Pi settings: disable changelog popup, quiet startup ---
PI_SETTINGS="${PI_DIR}/settings.json"
if [[ ! -f "${PI_SETTINGS}" ]]; then
    echo '{}' > "${PI_SETTINGS}"
fi
tmp=$(mktemp)
jq '. + {"quietStartup": true, "collapseChangelog": true}' "${PI_SETTINGS}" > "$tmp" && mv "$tmp" "${PI_SETTINGS}"

# --- Clean up legacy extension copy from /homeassistant volume ---
if [[ -d /homeassistant/.pi/extensions/home-assistant ]]; then
    bashio::log.info 'Removing legacy extension from /homeassistant/.pi/extensions/home-assistant...'
    rm -rf /homeassistant/.pi/extensions/home-assistant
fi

# --- Convenient symlinks ---
ln -sf /homeassistant /root/homeassistant
ln -sf /share /root/share

bashio::log.info 'Pi Agent environment ready.'
